# resources/templates/default/cape-web.service.erb
[Unit]
Description=CAPE WSGI app
Documentation=https://github.com/kevoreilly/CAPEv2
Wants=cape-rooter.service
After=cape-rooter.service

[Service]
Environment=VIRTUAL_ENV=/opt/redborder-cape/venv
Environment=PATH=/opt/redborder-cape/venv/bin:/usr/bin:/bin
WorkingDirectory=/opt/CAPEv2/web
ExecStart=/opt/redborder-cape/venv/bin/poetry run python manage.py runserver_plus <%= @cape_web_ip %>:<%= @cape_web_port %> --traceback --keep-meta-shutdown
User=cape
Group=cape

[Install]
WantedBy=multi-user.target
